# _backup_tools/com.claude-backup-notify.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.claude-backup-notify</string>

    <!-- UPDATE: Replace /Users/youruser with your actual home directory path -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/youruser/Projects/claude_backup/_backup_tools/claude-backup-notify.sh</string>
    </array>

    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>7</integer>
        <key>Minute</key>
        <integer>30</integer>
    </dict>

    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
    </dict>

    <!-- UPDATE: Replace /Users/youruser with your actual home directory path -->
    <key>StandardOutPath</key>
    <string>/Users/youruser/Library/Logs/claude-backup/launchd-notify.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/youruser/Library/Logs/claude-backup/launchd-notify.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
